# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:django]
command=poetry run python manage.py runserver 0.0.0.0:8000
directory=/src/app/uptime-monitor
autostart=true
autorestart=true
stderr_logfile=/var/log/django.err.log
stdout_logfile=/var/log/django.out.log
priority=100


[program:celery_worker]
command=poetry run celery -A uptimemonitor worker --loglevel=INFO
directory=/src/app/uptime-monitor
autostart=true
autorestart=true
stderr_logfile=/var/log/celery_worker.err.log
stdout_logfile=/var/log/celery_worker.out.log
priority=200

[program:celery_beat]
command=poetry run celery -A uptimemonitor beat -l info --scheduler django_celery_beat.schedulers:DatabaseScheduler
directory=/src/app/uptime-monitor
autostart=true
autorestart=true
stderr_logfile=/var/log/celery_beat.err.log
stdout_logfile=/var/log/celery_beat.out.log
priority=300
