# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:django]
command=poetry run python /src/app/uptime-monitor/manage.py runserver 0.0.0.0:8000
autostart=true
autorestart=true
